# com.fuggerbot.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fuggerbot.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ryanfisher/fuggerbot/run_monitor.py</string>
        <string>--interval</string>
        <string>60</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ryanfisher/fuggerbot</string>
    <key>StandardOutPath</key>
    <string>/Users/ryanfisher/fuggerbot/data/logs/monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ryanfisher/fuggerbot/data/logs/monitor.error.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>

